# z-nomp-util/util/remove_claymore_fees.crontab
*/15 * * * * ~/z-nomp/util/remove_claymore_fees.sh