# z-nomp-util/util/remove_claymore_fees.crontab
*/15 * * * * ~/rig-works/z-nomp-util/util/remove_claymore_fees.sh
